# ribxz_update.service
# /etc/systemd/system/ribxz-update.service
[Unit]
Description=RibXZ Database Update Service
After=docker.service
Requires=docker.service

[Service]
Type=oneshot
User=rtviii
Group=docker
WorkingDirectory=/home/rtviii/dev/riboxyz

# Environment setup
Environment=PYTHONUNBUFFERED=1
EnvironmentFile=/home/rtviii/dev/riboxyz/.env

# Execution
ExecStart=/bin/bash -c '\
    source /home/rtviii/dev/riboxyz/ribxzvenv/bin/activate && \
    python3 ribctl/ribd.py etl sync_all -t PTC -t STRUCTURE_PROFILE -t MMCIF --workers 6 --delay 2 && \
    python3 ribctl/ribd.py db upload_all -w 6'

# Logging
StandardOutput=append:/var/log/ribxz/update.log
StandardError=append:/var/log/ribxz/update.error.log

# Restart settings
Restart=no
TimeoutStartSec=0

# Security settings
NoNewPrivileges=yes
ProtectSystem=full
ReadWritePaths=/var/log/ribxz /home/rtviii/dev/riboxyz

[Install]
WantedBy=multi-user.target

# /etc/systemd/system/ribxz-update.timer
[Unit]
Description=Run RibXZ database update every two weeks
Requires=ribxz-update.service

[Timer]
OnCalendar=*-*-1,15 00:00:00
RandomizedDelaySec=3600
Persistent=true

[Install]
WantedBy=timers.target